# .config/systemd/user/blueman.service
[Unit]
Description=Bluetooth Manager

[Service]
Type=idle
ExecStart=/usr/bin/blueman-applet

[Install]
WantedBy=Xorg.target
