# .config/systemd/user/fcitx.service
[Unit]
Description=Input method framework

[Service]
Type=forking
ExecStart=/usr/bin/fcitx5 -d

[Install]
WantedBy=Xorg.target
