# .config/systemd/user/fcitx5.service
[Unit]
Description=fcitx5

[Service]
Type=forking
ExecStart=/usr/bin/fcitx5 -d

[Install]
WantedBy=X11.target
WantedBy=Wayland.target
