# .config/systemd/user/flameshot.service
[Unit]
Description=Screenshot software

[Service]
Type=idle
ExecStart=/usr/bin/flameshot

[Install]
WantedBy=Xorg.target
