# .config/systemd/user/kdeconnect.service
[Unit]
Description=kdeconnect

[Service]
ExecStart=/usr/bin/kdeconnect-indicator

[Install]
WantedBy=X11.target
WantedBy=Wayland.target
