# .config/systemd/user/network-manager-applet.service
[Unit]
Description=network-manager-applet

[Service]
ExecStart=/usr/bin/nm-applet

[Install]
WantedBy=X11.target
WantedBy=Wayland.target
