# .config/systemd/user/picom.service
[Unit]
Description=X compositor

[Service]
Type=idle
Environment="LANG=zh_CN.UTF-8"
ExecStart=/usr/bin/picom

[Install]
WantedBy=Xorg.target
