# .config/systemd/user/policykit-agent.service
[Unit]
Description=LXQt policykit authentication agent

[Service]
Type=idle
ExecStart=/usr/bin/lxqt-policykit-agent

[Install]
WantedBy=Xorg.target
