# .config/systemd/user/udiskie.service
[Unit]
Description=udiskie

[Service]
ExecStart=/usr/bin/udiskie --smart-tray

[Install]
WantedBy=X11.target
WantedBy=Wayland.target
