# .config/systemd/user/Wayland.target
[Unit]
Description=Run desktop applications under Wayland
