# .config/systemd/user/Xorg.target
[Unit]
Description=Launch desktop applications
Requires=default.target
After=default.target
