# .config/systemd/user/xsettingsd.service
[Unit]
Description=xsettingsd

[Service]
ExecStart=/usr/bin/xsettingsd

[Install]
WantedBy=X11.target
