# .config/systemd/user/plasma-polkit-agent.service.d/override.conf
[Install]
WantedBy=X11.target
WantedBy=Wayland.target